# ops/cardputer-proxy/cardputer-proxy.service
[Unit]
Description=Cardputer-AI claude-proxy (OAI-compatible to Anthropic)
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=cardputer-proxy
Group=cardputer-proxy
WorkingDirectory=/opt/cardputer-proxy

# Secrets are read directly from /etc/cardputer-proxy/secrets/ at request
# time. The directory is 0750 root:cardputer-proxy and files are 0640
# root:cardputer-proxy, so the service user can read any secret_ref the
# profile catalog references without needing a unit-file edit per provider.
Environment=CARDPUTER_PROXY_SECRETS_DIR=/etc/cardputer-proxy/secrets

EnvironmentFile=-/etc/cardputer-proxy/cardputer-proxy.env
ExecStart=/opt/cardputer-proxy/.venv/bin/python -m cardputer_proxy

Restart=on-failure
RestartSec=2s

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/cardputer-proxy
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
